# motdgen.service
[Unit]
Description=motdgen.service

[Service]
ExecStart=/usr/bin/motdgen
Type=oneshot
RemainAfterExit=yes
